# systemd/voicenote-processor.service
[Unit]
Description=Process a new voice note via transcription and Fabric
Documentation=https://github.com/mishrasidhant/voicenote-automation

[Service]
Type=oneshot
# The placeholder __INSTALL_DIR__ will be replaced by the install script.
ExecStart=__INSTALL_DIR__/bin/process-voicenote.sh %f
