# systemd/voicenote-processor.path
[Unit]
Description=Watch for new voice notes to process
Documentation=https://github.com/mishrasidhant/voicenote-automation

[Path]
PathChanged=%h/voicenotes/in/

[Install]
WantedBy=default.target
